# smartdns-proxy.service
[Unit]
Description=SmartDNS Proxy Server
After=network.target
Documentation=https://github.com/yourusername/smartdns-proxy

[Service]
Type=simple
User=root
WorkingDirectory=/home/smartdns
ExecStart=/home/smartdns/smartdns-proxy
Restart=on-failure
RestartSec=5s

# Security
PrivateTmp=true
NoNewPrivileges=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smartdns-proxy

[Install]
WantedBy=multi-user.target
